# etc/systemd/system/entralinked.service
[Unit]
Description=Entralinked Java Service
After=network.target

[Service]
# Ensure the service runs from the correct directory
WorkingDirectory=/root/entralinked/
# Use the full path to the java executable
ExecStart=/usr/bin/java -jar /root/entralinked/entralinked.jar disablegui
# Restart automatically if the process crashes
Restart=always
RestartSec=5
# Treat exit code 143 (SIGTERM) as a clean exit
SuccessExitStatus=143
# Standard output/error will go to the system journal
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target